# build-install-scripts/system-service/template.service
[Unit]
Description = description
After = network.target

[Service]
ExecStart = /root/root.sh

[Install]
WantedBy = multi-user.target
